# GitPage/dist/etc/systemd/kexec-load.service
# Automatically generated by systemd-sysv-generator

[Unit]
Description=LSB: Load kernel image with kexec
Documentation=man:kexec(8)
DefaultDependencies=no
Before=shutdown.target umount.target final.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/kexec -l /vmlinuz --initrd=/initrd.img --reuse-cmdline

[Install]
WantedBy=kexec.target
